# deploy_tools/gunicorn-upstart.template.conf
description "gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid markg
chdir /home/markg/Sites/SITENAME/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket GSD.wsgi:application
